# etc/NetworkManager/dispatcher.d/10-i3blocks.sh
#!/bin/sh

pkill -SIGRTMIN+9 i3blocks

exit 0
